# data/miatoll-settings.service.in
[Unit]
Description=Miatoll Settings

[Service]
Type=exec
ExecStart=@BIN_DIR@/miatoll-settings

[Install]
WantedBy=gnome-session-initialized.target
