# data/mobile-settings.service.in
[Unit]
Description=Mobile Settings

[Service]
Type=exec
ExecStart=@BIN_DIR@/mobile-settings

[Install]
WantedBy=graphical-session.target
